# themeswitcher.timer
[Unit]
Description=Run themeswitcher shecduly

[Timer]
OnUnitActiveSec=1h
Unit=themeswitcher.service

[Install]
WantedBy=graphical.target
